# roles/xray/templates/install_autoupdate.timer.j2
[Unit]
Description=Xray auto-update timer

[Timer]
OnCalendar={{ xray_autoupdate_calendar }}
Persistent=true
RandomizedDelaySec={{ xray_autoupdate_random_delay }}

[Install]
WantedBy=timers.target
